# bin for ubuntu 20.04/live555proxy.service
[Unit]
Description = do some RTSP proxying
After = network.target

[Service]
Type=forking
User=root
Restart=always
ExecStart=live555ProxyServer -p 8554 rtsp://USERNAME:PASSWORD@CAM_IP:554/h264Preview_01_sub

[Install]
WantedBy=multi-user.target
